# packages/zramswap/files/zram-swap.service
[Unit]
Description=Service enabling compressing RAM with zRam
After=swap.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/zramswapon
ExecStop=/usr/local/sbin/zramswapoff
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

